# etc/init.d/mac-filter
#!/bin/sh /etc/rc.common

START=99
USE_PROCD=1

NFT_TABLE="mac_filter"
NFT_CHAIN="block_mac"

service_triggers() {
    procd_add_reload_trigger "mac-filter"
}

reload_service() {
    # 清除旧规则
    nft delete table inet $NFT_TABLE 2>/dev/null
    
    # 创建新表和链
    nft create table inet $NFT_TABLE
    nft create chain inet $NFT_TABLE $NFT_CHAIN {
        type filter hook forward priority 0 \; policy accept \;
    }

    # 加载阻断规则
    local enabled=$(uci -q get mac-filter.@rule[0].enabled)
    local mac=$(uci -q get mac-filter.@rule[0].mac)
    
    if [ "$enabled" = "1" ] && [ -n "$mac" ]; then
        nft add rule inet $NFT_TABLE $NFT_CHAIN ether saddr $mac counter drop
    fi
}
