# etc/init.d/mac-filter
#!/bin/sh /etc/rc.common

START=99
USE_PROCD=1

RULES_COMMENT="MAC_FILTER"  # 规则标识

start_service() {
    /etc/init.d/mac-filter reload
}

flush_rules() {
    # 清除所有带指定注释的规则
    iptables-save | grep -v "$RULES_COMMENT" | iptables-restore
}

reload_service() {
    # 清除旧规则
    flush_rules

    # 从 UCI 配置加载新规则
    local blocked_mac=$(uci -q get mac-filter.@rule[0].mac)
    local enabled=$(uci -q get mac-filter.@rule[0].enabled)

    if [ "$enabled" = "1" ] && [ -n "$blocked_mac" ]; then
        # 在 FORWARD 链插入规则（根据网络拓扑可能需要调整链）
        iptables -I FORWARD -m mac --mac-source "$blocked_mac" -m comment --comment "$RULES_COMMENT" -j DROP
    fi
}
